# deploy/backup/hms-backup.service
[Unit]
Description=HMS daily backup (DB + uploads)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hms-backup.sh
